# BenutzerDB/debian/raumzeitbenutzerdb.cron.d
1 0 * * * root /usr/bin/raumzeitbenutzerdb-pinrevoker
