# BenutzerDB/debian/raumzeitbenutzerdb.cron.d
unset MAILTO
1 0 * * * /usr/bin/raumzeitbenutzerdb-pinrevoker
